# lib/systemd/system/onion-grater.service.d/30_whonix.conf
## Copyright (C) 2012 - 2022 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

#### meta start
#### project Whonix
#### category networking and tor-control
#### gateway_only yes
#### description
## onion-grater systemd unit file extension
## Necessary changes for remote clients
#### meta end

[Service]
## Start onion-grater with listening on interface eth1 explicitly.
## Allow multiple qualifiers to match on different profiles
Environment="MAIN=--listen-interface eth1 --multiple-matches"
## Set profiles directories
Environment="DIRS=--profile-dir /etc/onion-grater-merger.d --profile-dir /usr/local/etc/onion-grater-merger.d"
## Enable debug mode.
#Environment="DEV=--debug"
## Enable complain mode, debug and allows everything (insecure in production).
#Environment="DEV=--complain"
## Override the above variables via custom file, prefix with '-' to avoid
## problems with missing file
EnvironmentFile=-/etc/onion-grater-systemd.conf
EnvironmentFile=-/usr/local/etc/onion-grater-systemd.conf

## Clear previous start line
ExecStart=
## Use variables assigned to arguments
ExecStart=/usr/lib/onion-grater $MAIN $DIRS $DEV
